# lfs258/s_02/foo.slice
[Unit]
Description=stress slice
[Slice]
CPUQuota=30%
